# systemd/lab-brain-status.timer
[Unit]
Description=Run Lab Mind status check every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
